# .config/autostart/lxpanel.desktop
[Desktop Entry]
Type=Application
Name=Bottom Panel
TryExec=lxpanel
Exec=lxpanel
Hidden=true
