# .config/autostart/pidgin.desktop
[Desktop Entry]
Type=Application
Name=Pidgin messenger
TryExec=pidgin
Exec=pidgin
Hidden=false
