# .config/autostart/skype.desktop
[Desktop Entry]
Type=Application
Name=Skype messenger
TryExec=skype
Exec=skype
Hidden=false
