# .config/autostart/tilda.desktop
[Desktop Entry]
Type=Application Name=Drop-down console
TryExec=tilda
Exec=tilda
Hidden=true
Comment=

