# .config/autostart/tint2.desktop
[Desktop Entry]
Type=Application
Name=Bottom Panel
TryExec=tint2
Exec=tint2
Hidden=true
