# .config/autostart/yeahconsole.desktop
[Desktop Entry]
Type=Application
Name=Drop-down console
TryExec=yeahconsole
Exec=yeahconsole
Hidden=true
Comment=

